# nspanel-lovelace-ui/rootfs/etc/s6-overlay/s6-rc.d/mqtt-manager/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

## Run your program
exec /usr/bin/mqtt-manager/run.sh
